# config/launchd/com.wyattowalsh.mcphub.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.wyattowalsh.mcphub</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/ww/dev/projects/agents/scripts/mcphub/launch-agent-run.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/ww/dev/projects/agents</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/com.wyattowalsh.mcphub.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/com.wyattowalsh.mcphub.err.log</string>
</dict>
</plist>
